# ops/hexo-preview.service
[Unit]
Description=Hexo Preview Server
After=network.target

[Service]
Type=simple
User=lucas_wsl
WorkingDirectory=/home/lucas_wsl/workspace/blog/liulangrenaaa.github.io
ExecStart=/home/lucas_wsl/workspace/blog/liulangrenaaa.github.io/node_modules/.bin/hexo server --port 4000 --host 0.0.0.0
Restart=always
RestartSec=5
Environment=NODE_ENV=development

[Install]
WantedBy=multi-user.target
